# infra/services/starlight/starlight.socket
[Socket]
ListenStream=80
ListenStream=443

[Install]
WantedBy=sockets.target
